# variants/with-shibboleth/supervisor/shibd.conf
[program:shibd]
command=shibd -F
priority=900
username=www-data
autostart=true
autorestart=true