# variants/with-shibboleth/supervisor/shibd.conf
[program:shibd]
command=shibd -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfil=/dev/stderr
stderr_logfile_maxbytes=0
priority=900
username=www-data
autostart=true
autorestart=true
startsecs=30